# deployment/celery-worker-local.service
[Unit]
Description=TEKDB Celery Worker
After=network.target postgresql.service redis-server.service

[Service]
Type=simple
User=vagrant
Group=vagrant
WorkingDirectory=/usr/local/apps/TEKDB/TEKDB
Environment="DJANGO_SETTINGS_MODULE=TEKDB.settings"
ExecStart=/usr/local/apps/env/bin/celery -A TEKDB worker \
    --loglevel=info \
    --logfile=/var/log/celery/tekdb-worker.log
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target