# deployment/rc.local
#!/bin/sh -e

TIME=`/bin/date +20%y/%m/%d-%H:%M:%S`

/usr/sbin/service uwsgi start

/usr/local/apps/TEKDB/deployment/restart_nginx.sh &

exit 0
